# files/supervisord-elasticsearch.conf
[program:elasticsearch]
command=/elasticsearch/src/bin/elasticsearch -f -Des.config=/elasticsearch/etc/elasticsearch.yml
autostart=true
autorestart=true

